# etc/systemd/system/scanbd.service.d/service.conf
[Service]
RestartPreventExitStatus=11
TimeoutStopSec=10min
KillMode=process
